# root/etc/services.d/seafile-client/run
#!/usr/bin/with-contenv bash

exec /usr/bin/seaf-cli start
exec /usr/bin/data
